# infra/model-services/systemd-user/qwen36-vllm.service
[Unit]
Description=Qwen3.6 35B FP8 vLLM API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/maoyd
ExecStartPre=/usr/bin/test -f /home/maoyd/models/Qwen3.6-35B-A3B-FP8-modelscope/config.json
ExecStart=/usr/bin/bash /home/maoyd/modles_setup/serve_qwen36_fp8_vllm_newenv.sh
Restart=no
RestartSec=10
TimeoutStartSec=900
TimeoutStopSec=120
SuccessExitStatus=0 143
KillMode=control-group
SyslogIdentifier=qwen36-vllm
Environment=HOME=/home/maoyd
Environment=TMPDIR=/tmp
Environment=PATH=/home/maoyd/.local/bin:/home/maoyd/bin:/home/maoyd/miniconda3/bin:/usr/local/bin:/usr/bin:/bin
Environment=ENV_DIR=/home/maoyd/miniconda3/envs/qwen36_fp8_vllm
Environment=QWEN36_VLLM_BIN=/home/maoyd/miniconda3/envs/qwen36_fp8_vllm/bin/vllm
Environment=LANGUAGE_MODEL_ONLY=0
Environment=VLLM_NO_USAGE_STATS=1
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
