# infra/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:backend]
command=/bin/sh -lc 'exec "${SIQ_BACKEND_ROOT:-%(here)s/../../apps/api}/start.sh"'
directory=%(here)s/../../apps/api
autostart=true
autorestart=true
stderr_logfile=%(here)s/../../var/logs/backend.err.log
stdout_logfile=%(here)s/../../var/logs/backend.out.log
environment=SIQ_PROJECT_ROOT="%(here)s/../..",SIQ_DATA_ROOT="%(here)s/../../data",SIQ_RUNTIME_ROOT="%(here)s/../../var",SIQ_ARTIFACTS_ROOT="%(here)s/../../artifacts",SIQ_DATASETS_ROOT="%(here)s/../../datasets"

[program:frontend]
command=npm run dev -- --host 0.0.0.0 --port 15173
directory=%(here)s/../../apps/web
autostart=true
autorestart=true
stderr_logfile=%(here)s/../../var/logs/frontend.err.log
stdout_logfile=%(here)s/../../var/logs/frontend.out.log

[program:report-finder]
command=/bin/sh -lc 'cd "${SIQ_REPORT_FINDER_ROOT:-%(here)s/../../services/market-report-finder}" && MARKET_REPORT_DOWNLOAD_DIR="${SIQ_REPORT_DOWNLOADS_ROOT:-${SIQ_DATA_ROOT:-%(here)s/../../data}/market-report-finder/downloads}" exec uv run python -m uvicorn market_report_finder_service.app:app --host 127.0.0.1 --port 18000'
directory=%(here)s/../..
autostart=true
autorestart=true
stderr_logfile=%(here)s/../../var/logs/report-finder.err.log
stdout_logfile=%(here)s/../../var/logs/report-finder.out.log
environment=SIQ_PROJECT_ROOT="%(here)s/../..",SIQ_DATA_ROOT="%(here)s/../../data",SIQ_RUNTIME_ROOT="%(here)s/../../var",SIQ_ARTIFACTS_ROOT="%(here)s/../../artifacts",SIQ_DATASETS_ROOT="%(here)s/../../datasets"

[program:pdf-parser]
command=%(here)s/../../apps/pdf-parser/run.sh
directory=%(here)s/../../apps/pdf-parser
autostart=true
autorestart=true
environment=HOST="127.0.0.1",PORT="15000",SIQ_PROJECT_ROOT="%(here)s/../..",SIQ_DATA_ROOT="%(here)s/../../data",SIQ_RUNTIME_ROOT="%(here)s/../../var",SIQ_ARTIFACTS_ROOT="%(here)s/../../artifacts",SIQ_DATASETS_ROOT="%(here)s/../../datasets",SIQ_PDF2MD_ROOT="%(here)s/../../apps/pdf-parser",SIQ_PDF2MD_DATA_DIR="%(here)s/../../data/pdf-parser"
stderr_logfile=%(here)s/../../var/logs/pdf-parser.err.log
stdout_logfile=%(here)s/../../var/logs/pdf-parser.out.log

[program:document-parser]
command=%(here)s/../../apps/document-parser/run.sh
directory=%(here)s/../../apps/document-parser
autostart=true
autorestart=true
environment=HOST="127.0.0.1",PORT="15010",SIQ_PROJECT_ROOT="%(here)s/../..",SIQ_DATA_ROOT="%(here)s/../../data",SIQ_RUNTIME_ROOT="%(here)s/../../var",SIQ_ARTIFACTS_ROOT="%(here)s/../../artifacts",SIQ_DATASETS_ROOT="%(here)s/../../datasets",SIQ_DOCUMENT_PARSER_ROOT="%(here)s/../../apps/document-parser",SIQ_DOCUMENT_PARSE_DATA_DIR="%(here)s/../../data/document-parser"
stderr_logfile=%(here)s/../../var/logs/document-parser.err.log
stdout_logfile=%(here)s/../../var/logs/document-parser.out.log

[group:hermes_agents]
programs=hermes_assistant,hermes_analysis,hermes_factchecker,hermes_tracking,hermes_legal

[program:hermes_assistant]
command=/bin/sh -lc 'exec "%(here)s/../../scripts/hermes/run_gateway.sh" siq_assistant'
directory=%(here)s/../..
autostart=true
autorestart=true
stderr_logfile=%(here)s/../../var/logs/hermes_assistant.err.log
stdout_logfile=%(here)s/../../var/logs/hermes_assistant.out.log

[program:hermes_analysis]
command=/bin/sh -lc 'exec "%(here)s/../../scripts/hermes/run_gateway.sh" siq_analysis'
directory=%(here)s/../..
autostart=true
autorestart=true
stderr_logfile=%(here)s/../../var/logs/hermes_analysis.err.log
stdout_logfile=%(here)s/../../var/logs/hermes_analysis.out.log

[program:hermes_factchecker]
command=/bin/sh -lc 'exec "%(here)s/../../scripts/hermes/run_gateway.sh" siq_factchecker'
directory=%(here)s/../..
autostart=true
autorestart=true
stderr_logfile=%(here)s/../../var/logs/hermes_factchecker.err.log
stdout_logfile=%(here)s/../../var/logs/hermes_factchecker.out.log

[program:hermes_tracking]
command=/bin/sh -lc 'exec "%(here)s/../../scripts/hermes/run_gateway.sh" siq_tracking'
directory=%(here)s/../..
autostart=true
autorestart=true
stderr_logfile=%(here)s/../../var/logs/hermes_tracking.err.log
stdout_logfile=%(here)s/../../var/logs/hermes_tracking.out.log

[program:hermes_legal]
command=/bin/sh -lc 'exec "%(here)s/../../scripts/hermes/run_gateway.sh" siq_legal'
directory=%(here)s/../..
autostart=true
autorestart=true
stderr_logfile=%(here)s/../../var/logs/hermes_legal.err.log
stdout_logfile=%(here)s/../../var/logs/hermes_legal.out.log
