# infra/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:backend]
command=/bin/sh -lc 'exec "${SIQ_BACKEND_ROOT:-%(here)s/../../apps/api}/start.sh"'
directory=%(here)s/../../apps/api
autostart=true
autorestart=true
stderr_logfile=%(here)s/../../data/logs/backend.err.log
stdout_logfile=%(here)s/../../data/logs/backend.out.log
environment=SIQ_PROJECT_ROOT="%(here)s/../.."

[program:frontend]
command=npm run dev -- --host 0.0.0.0 --port 15173
directory=%(here)s/../../apps/web
autostart=true
autorestart=true
stderr_logfile=%(here)s/../../data/logs/frontend.err.log
stdout_logfile=%(here)s/../../data/logs/frontend.out.log

[program:report-finder]
command=/bin/sh -lc 'cd "${SIQ_REPORT_FINDER_ROOT:-%(here)s/../../services/report-finder}" && exec uv run python -m uvicorn report_finder_service.app:app --host 127.0.0.1 --port 18000'
directory=%(here)s/../..
autostart=true
autorestart=true
stderr_logfile=%(here)s/../../data/logs/report-finder.err.log
stdout_logfile=%(here)s/../../data/logs/report-finder.out.log

[program:pdf-parser]
command=%(here)s/../../apps/pdf-parser/run.sh
directory=%(here)s/../../apps/pdf-parser
autostart=true
autorestart=true
environment=HOST="127.0.0.1",PORT="15000",SIQ_PROJECT_ROOT="%(here)s/../..",SIQ_PDF2MD_ROOT="%(here)s/../../apps/pdf-parser",SIQ_PDF2MD_DATA_DIR="%(here)s/../../data/pdf-parser"
stderr_logfile=%(here)s/../../data/logs/pdf-parser.err.log
stdout_logfile=%(here)s/../../data/logs/pdf-parser.out.log

[group:hermes_agents]
programs=hermes_assistant,hermes_analysis,hermes_factchecker,hermes_tracking,hermes_legal

[program:hermes_assistant]
command=/bin/sh -lc 'profile_dir="$("%(here)s/../../scripts/hermes/profile_dir.sh" siq_assistant)" && cd "$profile_dir" && export HERMES_HOME="$profile_dir" && exec hermes gateway run --replace --accept-hooks'
directory=%(here)s/../..
autostart=true
autorestart=true
stderr_logfile=%(here)s/../../data/logs/hermes_assistant.err.log
stdout_logfile=%(here)s/../../data/logs/hermes_assistant.out.log

[program:hermes_analysis]
command=/bin/sh -lc 'profile_dir="$("%(here)s/../../scripts/hermes/profile_dir.sh" siq_analysis)" && cd "$profile_dir" && export HERMES_HOME="$profile_dir" && exec hermes gateway run --replace --accept-hooks'
directory=%(here)s/../..
autostart=true
autorestart=true
stderr_logfile=%(here)s/../../data/logs/hermes_analysis.err.log
stdout_logfile=%(here)s/../../data/logs/hermes_analysis.out.log

[program:hermes_factchecker]
command=/bin/sh -lc 'profile_dir="$("%(here)s/../../scripts/hermes/profile_dir.sh" siq_factchecker)" && cd "$profile_dir" && export HERMES_HOME="$profile_dir" && exec hermes gateway run --replace --accept-hooks'
directory=%(here)s/../..
autostart=true
autorestart=true
stderr_logfile=%(here)s/../../data/logs/hermes_factchecker.err.log
stdout_logfile=%(here)s/../../data/logs/hermes_factchecker.out.log

[program:hermes_tracking]
command=/bin/sh -lc 'profile_dir="$("%(here)s/../../scripts/hermes/profile_dir.sh" siq_tracking)" && cd "$profile_dir" && export HERMES_HOME="$profile_dir" && exec hermes gateway run --replace --accept-hooks'
directory=%(here)s/../..
autostart=true
autorestart=true
stderr_logfile=%(here)s/../../data/logs/hermes_tracking.err.log
stdout_logfile=%(here)s/../../data/logs/hermes_tracking.out.log

[program:hermes_legal]
command=/bin/sh -lc 'profile_dir="$("%(here)s/../../scripts/hermes/profile_dir.sh" siq_legal)" && cd "$profile_dir" && export HERMES_HOME="$profile_dir" && exec hermes gateway run --replace --accept-hooks'
directory=%(here)s/../..
autostart=true
autorestart=true
stderr_logfile=%(here)s/../../data/logs/hermes_legal.err.log
stdout_logfile=%(here)s/../../data/logs/hermes_legal.out.log
